# scripts/com.astralmaris.jikimi-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.astralmaris.jikimi-watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/proteus/astralmaris/ming-qiao/jikimi/scripts/jikimi-watchdog.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/Users/proteus/astralmaris/ming-qiao/jikimi/health/jikimi-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/proteus/astralmaris/ming-qiao/jikimi/health/jikimi-watchdog.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/proteus</string>
    </dict>
</dict>
</plist>
